# dist/systemd/netfyr.socket
[Unit]
Description=Netfyr Varlink socket

[Socket]
ListenStream=/run/netfyr/netfyr.sock
SocketMode=0666

[Install]
WantedBy=sockets.target
